# rootfs/lib/systemd/system/update-dns.service
[Unit]
Description=Update DNS record on boot
After=network.target

[Service]
Type=simple
User=nobody
EnvironmentFile=-/etc/jupyterhub/jupyterhub.conf
EnvironmentFile=/etc/environment
ExecStart=/usr/local/bin/ec2update_dns "${DOMAIN}"
Restart=no

[Install]
WantedBy=multi-user.target