# deb/postinst
#!/bin/bash

echo "Enabling DEA init service"
systemctl enable dea-init.service dea-early-init.service
